# i2pd-docker-registry/supervisor-registry.conf
[program:docker-registry]
command="service docker-registry start"

