# i2pd-docker-hostbase/supervisor-registry.conf
[program:gitweb]
command="service docker-registry start"

